# service/pdns_exporter/run
#!/bin/sh

EXPORTER_LISTEN_ADDRESS=${EXPORTER_LISTEN_ADDRESS:-:9120}
EXPORTER_TELEMETRY_PATH=${EXPORTER_TELEMETRY_PATH:-/metrics}

sleep 5
exec /pdns_exporter \
  -web.listen-address=${EXPORTER_LISTEN_ADDRESS} \
  -web.telemetry-path=${EXPORTER_TELEMETRY_PATH}